# deploy/systemd/iona-rpc.service
[Unit]
Description=IONA RPC Node (Public Endpoint)
After=network-online.target
Wants=network-online.target
# Start after producers are likely up
After=iona-val2.service iona-val3.service iona-val4.service

[Service]
Type=simple
User=iona
Group=iona
ExecStart=/usr/local/bin/iona-node --config /etc/iona/rpc.toml
Restart=on-failure
RestartSec=2
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/iona/rpc
LimitNOFILE=1048576

# Environment
Environment=RUST_LOG=info
EnvironmentFile=-/etc/iona/rpc.env

[Install]
WantedBy=multi-user.target
